# playbooks/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/gunicorn -c /srv/webapps/altx-administration/src/gunicorn_config.py --reload altx-administration:app
directory={{ webapps_dir }}/{{ app_name }}/src
autostart=true
autorestart=true
